# fatigue-sentinel/fatigue-sentinel.timer
[Unit]
Description=Daily 08:30 creative fatigue sentinel

[Timer]
OnCalendar=*-*-* 08:30:00
Persistent=true

[Install]
WantedBy=timers.target
