# night-queue/systemd/night-queue-resume.service
[Unit]
Description=Relaunch night-queue driver if pending tasks and no driver running
[Service]
Type=oneshot
ExecStart=%h/systems/night-queue/resume.sh
