# systemd/agents-weekly.service
[Unit]
Description=agents-weekly (migrated cron)
[Service]
Type=oneshot
ExecStart=%h/sha-systems/agents/run_agents.sh weekly
