# systemd/box-watchdog.service
[Unit]
Description=box-watchdog
[Service]
Type=oneshot
ExecStart=%h/systems/watchdog/box-watchdog.sh
