# systemd/box-doctor.service
[Unit]
Description=box-doctor daily self-diagnostic
[Service]
Type=oneshot
ExecStart=%h/systems/box-doctor/doctor.sh
