# systemd/box-doctor.timer
[Unit]
Description=box-doctor timer
[Timer]
OnCalendar=*-*-* 08:45:00
Persistent=true
[Install]
WantedBy=timers.target
