# systemd/bq-clickup-perf.service
[Unit]
Description=BQ -> ClickUp perf writeback

[Service]
Type=oneshot
ExecStart=%h/systems/bq-clickup-perf/run_perf_writeback.sh
