# systemd/comments-digest.service
[Unit]
Description=comments-digest (migrated cron)
[Service]
Type=oneshot
ExecStart=%h/sha-systems/comments-digest/run_digest.sh
