# systemd/creative-feedback-synth.service
[Unit]
Description=creative-feedback-synth (migrated cron)
[Service]
Type=oneshot
ExecStart=%h/systems/creative-feedback/run_feedback_synth.sh
