# systemd/gbrain-weekly.service
[Unit]
Description=gbrain-weekly (migrated cron)
[Service]
Type=oneshot
ExecStart=%h/sha-systems/gbrain-weekly/run_gbrain_weekly.sh
