# systemd/qa-smoke.service
[Unit]
Description=weekly QA smoke test of the ~/systems cron/agent estate
[Service]
Type=oneshot
ExecStart=%h/systems/qa/smoke.sh
