# systemd/research-deepdive.service
[Unit]
Description=Research lanes engine
[Service]
Type=oneshot
ExecStart=%h/systems/research-agent/bin/run-lanes.sh
