# systemd/research-deepdive.service
[Unit]
Description=research-deepdive (migrated cron)
[Service]
Type=oneshot
ExecStart=%h/systems/research-agent/bin/deepdive.sh
