# systemd/research-monitor.service
[Unit]
Description=research-monitor (migrated cron)
[Service]
Type=oneshot
ExecStart=%h/systems/research-agent/bin/monitor.sh
