# systemd/sha-weekly-report.service
[Unit]
Description=sha-weekly-report (migrated cron)
[Service]
Type=oneshot
ExecStart=%h/sha-systems/sha-weekly-report/run_report.sh
