# systemd/tablet-screen-off.service
[Unit]
Description=Dashboard tablet screen off (overnight OLED/battery saver)
After=network-online.target

[Service]
Type=oneshot
User=tomas
Environment=ADB=/usr/bin/adb
Environment=TABLET_IP=192.168.0.160
Environment=HOME=/home/tomas
ExecStart=/home/tomas/systems/tablet-screen.sh off
