# systemd/tablet-screen-on.timer
[Unit]
Description=Tablet screen ON at 07:00 daily

[Timer]
OnCalendar=*-*-* 07:00:00
Persistent=false

[Install]
WantedBy=timers.target
