# systemd/task-lessons-synth.service
[Unit]
Description=task-lessons-synth (cluster lessons -> gated promotion proposals)
[Service]
Type=oneshot
ExecStart=%h/systems/task-lessons/run_lessons_synth.sh
