# systemd/task-lessons-synth.timer
[Unit]
Description=task-lessons-synth timer
[Timer]
OnCalendar=Tue *-*-* 05:30:00
Persistent=true
[Install]
WantedBy=timers.target
