# systemd/winners-refresh.service
[Unit]
Description=winners-refresh (migrated cron)
[Service]
Type=oneshot
ExecStart=%h/systems/winners-refresh/run_winners_refresh.sh
